# .docker/services/roadrunner.conf
[program:roadrunner]
command=/usr/local/bin/rr serve -d -c /etc/roadrunner/.rr.yaml
autorestart=true
priority=1
user=root
